# sv/nginx/run
#!/bin/bash

source /etc/envvars

mkdir -p /var/cache/nginx && chown nobody /var/cache/nginx

exec 2>&1
exec /usr/sbin/nginx -g "daemon off;"
