# sv/nginx/run
#!/bin/bash

source /etc/envvars

mkdir -p /var/nginx/cache && chown nobody /var/nginx/cache

exec 2>&1
exec /usr/sbin/nginx -g "daemon off;"
