# sv/redis/run
#!/bin/sh
exec 2>&1
exec redis-server /etc/redis.conf
