# sv/cron/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

crontab /crontab

exec 2>&1
exec cron -f
